# ops/macos/com.statarb-gemini.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.statarb-gemini.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_WITH_REPOSITORY_PATH/.venv/bin/trading-daemon</string>
    <string>--config</string>
    <string>REPLACE_WITH_REPOSITORY_PATH/config/production.json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_REPOSITORY_PATH</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>NetworkState</key>
    <true/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>REPLACE_WITH_REPOSITORY_PATH/var/log/daemon.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_REPOSITORY_PATH/var/log/daemon.stderr.log</string>
</dict>
</plist>
